# scripts/stcollectd.slc4
#!/bin/sh
#
# chkconfig: 2345 56 01
# description: Starts and stops the sFlow collector daemon
#
# pidfile: /var/run/$DAEMON.pid
# config:  /etc/$DAEMON.conf

ARGS=" -f 10 -p 6343 -i 10.146.45.25 -o /tmp/sflow -v -b 10 -s 2500"
DAEMON=stcollectd
DAEMONXML=stxmlrpcd
ARGSXML=""
DAEMON_HOME=/usr/local/bin
SERVICE="sFlow collecting"
SERVICEXML="sFlow XMLRPC"

# Source function library.
if [ -f /etc/init.d/functions ] ; then
  . /etc/init.d/functions
elif [ -f /etc/rc.d/init.d/functions ] ; then
  . /etc/rc.d/init.d/functions
else
  exit 1
fi

# Avoid using root's TMPDIR
unset TMPDIR

[ -f /etc/$DAEMON.conf ] || exit 6

RETVAL=0
RETVALXML=0

start() {
	echo -n $"Starting $SERVICE services: "
	daemon /usr/local/bin/$DAEMON $ARGS
	RETVAL=$?
	echo
	[ $RETVAL -eq 0 ] && touch /var/lock/subsys/$DAEMON || RETVAL=1
	echo -n $"Starting $SERVICEXML services: "
	daemon /usr/local/bin/$DAEMONXML $ARGSXML
	RETVALXML=$?
	echo
	[ $RETVALXML -eq 0 ] && touch /var/lock/subsys/$DAEMONXML || RETVALXML=1
	return  $RETVAL
}	

stop() {
	echo -n $"Shutting down $SERVICE services: "
	killproc $DAEMON
	RETVAL=$?
	echo
	[ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/$DAEMON

	echo -n $"Shutting down $SERVICEXML services: "
	killproc $DAEMONXML
	RETVALXML=$?
	echo
	[ $RETVALXML -eq 0 ] && rm -f /var/lock/subsys/$DAEMONXML
	return $RETVAL
}	

restart() {
	stop
	start
}	

rhstatus() {
	status $DAEMON
	RETVAL=$?
#	status $DAEMONXML

	# Handle this ourselves since something obviously doesnt like python daemons
<<<<<<< HEAD:scripts/stcollectd.slc4
        if [ -f /var/run/$DAEMONXML.pid -a -f /var/lock/subsys/$DAEMONXML ] ; then
                read pid < /var/run/$DAEMONXML.pid
		echo "$DAEMONXML (pid $pid) is running..."
        else
=======
	if [ -f /var/run/$DAEMONXML.pid -a -f /var/lock/subsys/$DAEMONXML ] ; then
		read pid < /var/run/$DAEMONXML.pid
		echo "$DAEMONXML (pid $pid) is running..."
	else
>>>>>>> b0cc51c... Updated SLC4 init script to show status on stxmlrpcd correctly:scripts/stcollectd.slc4
		echo "$DAEMONXML is stopped"
	fi

	RETVALXML=$?
	if [ $RETVAL -ne 0 ]; then
		return $RETVAL
	fi
	if [ $RETVALXML -ne 0 ]; then
		return $RETVALXML
	fi
}	

# Allow status as non-root.
if [ "$1" = status ]; then
       rhstatus
       exit $?
fi

case "$1" in
  start)
  	start
	;;
  stop)
  	stop
	;;
  restart)
  	restart
	;;
  status)
  	rhstatus
	;;
  *)
	echo $"Usage: $0 {start|stop|restart|status}"
	exit 2
esac

exit $?
